# ubuntu/docker/etc/udev/rules.d/local.rules

ACTION=="add", SUBSYSTEM=="module", KERNEL=="nf_conntrack", \
    RUN+="/lib/systemd/systemd-sysctl --prefix=/net/netfilter"
